# contrib/init/xcbcd.conf
description "XCBC Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env XCBCD_BIN="/usr/bin/xcbcd"
env XCBCD_USER="xcbc"
env XCBCD_GROUP="xcbc"
env XCBCD_PIDDIR="/var/run/xcbcd"
# upstart can't handle variables constructed with other variables
env XCBCD_PIDFILE="/var/run/xcbcd/xcbcd.pid"
env XCBCD_CONFIGFILE="/etc/xcbc/xcbc.conf"
env XCBCD_DATADIR="/var/lib/xcbcd"

expect fork

respawn
respawn limit 5 120
kill timeout 600

pre-start script
    # this will catch non-existent config files
    # xcbcd will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$XCBCD_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run xcbcd."
        echo "The setting must appear in $XCBCD_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following "
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"XCBC Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$XCBCD_PIDDIR"
    chmod 0755 "$XCBCD_PIDDIR"
    chown $XCBCD_USER:$XCBCD_GROUP "$XCBCD_PIDDIR"
    chown $XCBCD_USER:$XCBCD_GROUP "$XCBCD_CONFIGFILE"
    chmod 0660 "$XCBCD_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$XCBCD_PIDFILE" \
    --chuid $XCBCD_USER:$XCBCD_GROUP \
    --exec "$XCBCD_BIN" \
    -- \
    -pid="$XCBCD_PIDFILE" \
    -conf="$XCBCD_CONFIGFILE" \
    -datadir="$XCBCD_DATADIR" \
    -disablewallet \
    -daemon

